# package/pfSense-pkg-layer7/files/usr/local/etc/rc.d/layer7d
#!/bin/sh
#
# PROVIDES: layer7d
# REQUIRE: DAEMON NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="layer7d"
rcvar="${name}_enable"
load_rc_config "${name}"

: "${layer7d_enable:=NO}"

pidfile="/var/run/${name}.pid"
layer7d_bin="/usr/local/sbin/layer7d"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"
extra_commands="reload"
reload_cmd="${name}_reload"

layer7d_ensure_pf_tables()
{
	if [ -x /usr/local/libexec/layer7-pfctl ]; then
		/usr/local/libexec/layer7-pfctl ensure
		return $?
	fi
	return 0
}

# Le ${pidfile} e define _p como PID numerico, ou vazio se invalido/ausente.
# Alinhado a update-blacklists.sh / layer7-stats-collect.sh (F4.1).
layer7d_pid_from_file()
{
	_p=""
	[ -f "${pidfile}" ] || return 1
	if ! read -r _p <"${pidfile}" 2>/dev/null; then
		_p=""
		return 1
	fi
	_p=$(printf '%s' "${_p}" | sed 's/^[[:space:]]*//;s/[[:space:]]*$//')
	case "${_p}" in
	''|*[!0-9]*)
		_p=""
		return 1
		;;
	esac
	return 0
}

layer7d_start()
{
	local _ensure_rc

	/bin/mkdir -p /var/run
	if [ -f "${pidfile}" ]; then
		if layer7d_pid_from_file && /bin/kill -0 "${_p}" 2>/dev/null; then
			return 0
		fi
		/bin/rm -f "${pidfile}"
	fi
	layer7d_ensure_pf_tables
	_ensure_rc=$?
	if [ "${_ensure_rc}" -ne 0 ]; then
		echo "${name}: PF ensure reported errors (rc=${_ensure_rc}), starting daemon in degraded mode." >&2
	fi
	/usr/sbin/daemon -f -p "${pidfile}" "${layer7d_bin}"
	# `daemon` cria o pidfile tipicamente como 0600 root:wheel; operadores sem
	# privilegio nao conseguiam ler o PID e `service layer7d status` falhava
	# indevidamente. 0644 expoe apenas o numero do processo.
	/bin/chmod 0644 "${pidfile}" 2>/dev/null || true
}

layer7d_stop()
{
	local _i

	_p=""
	if [ -f "${pidfile}" ] && layer7d_pid_from_file; then
		:
	else
		_p=""
	fi
	/bin/rm -f "${pidfile}"
	# Sinal preferencial ao PID do pidfile (tipicamente o filho sob daemon(8)).
	[ -n "${_p}" ] && /bin/kill -TERM "${_p}" 2>/dev/null
	# Garantir paragem mesmo quando o pidfile aponta só para metade da arvore
	# (supervisor `daemon` + `layer7d`) ou ficou orfao — evita "Stop" no GUI
	# sem efeito e processos residuais.
	if [ -x /usr/bin/pkill ]; then
		/usr/bin/pkill -TERM -f "/usr/local/sbin/layer7d" 2>/dev/null || true
		_i=0
		while [ "${_i}" -lt 5 ]; do
			/usr/bin/pgrep -f "/usr/local/sbin/layer7d" >/dev/null 2>&1 || break
			_i=$(( _i + 1 ))
			sleep 1
		done
		/usr/bin/pkill -KILL -f "/usr/local/sbin/layer7d" 2>/dev/null || true
	fi
	return 0
}

layer7d_status()
{
	if [ -f "${pidfile}" ] && layer7d_pid_from_file && /bin/kill -0 "${_p}" 2>/dev/null; then
		echo "${name} is running as pid ${_p}."
		return 0
	fi
	echo "${name} is not running."
	return 1
}

layer7d_reload()
{
	if [ -f "${pidfile}" ] && layer7d_pid_from_file && /bin/kill -0 "${_p}" 2>/dev/null; then
		/bin/kill -HUP "${_p}"
		echo "${name}: reload signal sent."
		return 0
	fi
	echo "${name}: not running, starting instead."
	layer7d_start
}

run_rc_command "$1"
