# package/pfSense-pkg-layer7/pkg-plist
/usr/local/pkg/layer7.xml
/usr/local/pkg/layer7.inc
/etc/inc/priv/layer7.priv.inc
/usr/local/www/packages/layer7/layer7_status.php
/usr/local/www/packages/layer7/layer7_settings.php
/usr/local/www/packages/layer7/layer7_policies.php
/usr/local/www/packages/layer7/layer7_groups.php
/usr/local/www/packages/layer7/layer7_categories.php
/usr/local/www/packages/layer7/layer7_exceptions.php
/usr/local/www/packages/layer7/layer7_test.php
/usr/local/www/packages/layer7/layer7_events.php
/usr/local/www/packages/layer7/layer7_diagnostics.php
/usr/local/www/packages/layer7/layer7_blacklists.php
/usr/local/www/packages/layer7/layer7_bl_ajax.php
/usr/local/www/packages/layer7/layer7_reports.php
/usr/local/www/packages/layer7/layer7_reports_export.php
/usr/local/sbin/layer7d
/usr/local/etc/layer7.json.sample
/usr/local/etc/layer7-protos.txt.sample
/usr/local/etc/layer7/pf.conf.sample
/usr/local/etc/layer7/profiles.json
/usr/local/etc/rc.d/layer7d
/usr/local/libexec/layer7-pfctl
/usr/local/libexec/layer7-unbound-anti-doh
/usr/local/etc/layer7/update-blacklists.sh
/usr/local/etc/layer7/layer7-stats-collect.sh
/usr/local/etc/layer7/layer7-reports-collect.php
/usr/local/etc/layer7/layer7-stats-purge.sh
/usr/local/etc/layer7/lang/en.php
/usr/local/etc/layer7/lang/pt.php
/usr/local/share/pfSense-pkg-layer7/blacklists-signing-public-key.pem
/usr/local/share/pfSense-pkg-layer7/info.xml
@dir /etc/inc/priv
@dir /usr/local/etc/layer7/lang
@dir /usr/local/etc/layer7/reports
@dir /usr/local/etc/layer7/blacklists
@dir /usr/local/etc/layer7
@dir /usr/local/www/packages/layer7
